# content/etc/init/vnc-autores-controller.conf
### /etc/init/vnc-autores-controller.conf
# Controller service for /usr/local/sbin/vnc-autores
###
start on starting vnc-server or stopped vnc-server
task

export UPSTART_EVENTS
env AUTOMODE=0
env CONFIGFILE=/etc/default/vnc-autores

pre-start script
    [ ! -e $CONFIGFILE ] || . $CONFIGFILE

    [ $AUTOMODE -ne 1 ] || exit 0
    stop
end script

exec /usr/local/sbin/vnc-autores
