# build/content-tpl/etc/init/vnc-server.conf
start on stopped xbmc-loaded or started tty1 or vnc-server ACTION=START
stop on starting xbmc or vnc-server ACTION=STOP

env MULTI=1

# XBMC binary
env DAEMON="/usr/local/lib/kodi/kodi.bin"

nice +3

respawn

console none

pre-start exec modprobe -qb uinput

script
    if [ $(xbian-arch) = RPI ]; then
        exec /usr/local/sbin/rpi-vncserver
    else
        [ ! -e /etc/default/xbmc ] || . /etc/default/xbmc
        pgrep $(basename $DAEMON) >/dev/null && MULTI=2
        exec /usr/local/sbin/imx-vncserver -m $MULTI
    fi
end script

post-start script
    exit 0
    export oldmode=$(cat /sys/class/graphics/fb0/mode)
    while pgrep -x rpi-vncserver || pgrep -x imx-vncserver; do
        mode=$(cat /sys/class/graphics/fb0/mode)
        if [ "$mode" != "$oldmode" ]; then
            pkill -x imx-vncserver rpi-vncserver
            sleep 2
            break
        fi
        sleep 1
    done &
end script
